# supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
command=tor

[program:bitcoind]
command=bitcoind -debug=tor
